# dyno.cron
# Dyno!
* * * * * cd $DYNO_PATH; ./dyno.sh "$(date +%s)"
